# sys/fcitx.service
[Service]
ExecStart=fcitx
RestartSec=1s
Restart=always
