# sys/foot-server.service
[Service]
ExecStart=foot -s
RestartSec=1s
Restart=always
