# sys/radico8-player.service
[Service]
ExecStart=mpv http://radico8.xoc3.io &>/dev/null
Restart=always
RestartSec=1
