# sys/rivertile.service
[Service]
ExecStart=rivertile -view-padding 0 -outer-padding 0 -main-ratio .5
RestartSec=1s
Restart=always
