# sys/swayidle.service
[Service]
ExecStart=swayidle timeout 180 'waylock'
RestartSec=1s
Restart=always
