# sys/warpd.service
[Service]
ExecStart=%h/.local/bin/warpd -c %h/.dot/warpd-config -f --grid
Type=oneshot
Restart=no

[Install]
WantedBy=default.target
