# sys/warpd.service
[Service]
ExecStart=warpd -c /home/alan/.dot/warpd-config -f --grid
Type=oneshot
Restart=no

[Install]
WantedBy=default.target
