# sys/wl-gammarelay-rs.service
[Service]
ExecStart=wl-gammarelay-rs run
RestartSec=1s
Restart=always
